# ansible/roles/alpine/common/files/etc/supervisor.d/mailhog.supervisor.ini
[program:mailhog]
command=/opt/mailhog/mailhog